# packaging/debian/postremove.sh
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
	systemctl daemon-reload >/dev/null 2>&1 || true
fi

if [ "$1" = "purge" ]; then
	if command -v update-rc.d >/dev/null 2>&1; then
		update-rc.d upag remove >/dev/null 2>&1 || true
	fi
	rm -rf /var/lib/upag /run/upag
fi

exit 0
